# supervisord/supervisord-proxy.conf
[program:proxy.py]
command=python /proxy.py
numprocs=1
autostart=true
autorestart=true
