# supervisord/supervisord-openvpn.ini
[program:openvpn]
command=/start_vpn.sh
numprocs=1
autostart=true
autorestart=true
