# supervisord/supervisord-proxy.ini
[program:pysocks.py]
command=python /pysocks.py start
numprocs=1
autostart=true
autorestart=true
